# deploy/traceforge.service
[Unit]
Description=TraceForge evidence-bounded migration API
After=network-online.target
Wants=network-online.target

[Service]
User=traceforge
Group=traceforge
WorkingDirectory=/opt/traceforge
EnvironmentFile=/etc/traceforge/traceforge.env
ExecStart=/usr/local/bin/node /opt/traceforge/dist/server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ReadOnlyDirectories=/opt/traceforge
ReadWriteDirectories=/var/lib/traceforge
MemoryLimit=1536M
TasksMax=256
CPUQuota=200%
UMask=0027

[Install]
WantedBy=multi-user.target
